# autostart/ClipboardServer.desktop
[Desktop Entry]
Type=Application
Exec=io.github.q962.ClipboardServer
